# scripts/mammon-heartbeat.service
[Unit]
Description=MAMMON heartbeat dead-man check (oneshot, driven by mammon-heartbeat.timer)
Documentation=https://github.com/kpjmd/Mammon
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
# WorkingDirectory so the stdlib checker reads .env (ALERT_WEBHOOK,
# HEARTBEAT_FILE, HEARTBEAT_MAX_AGE_SECONDS) and data/heartbeat.json relative
# to the deploy dir. No poetry/venv needed — heartbeat_check.py is stdlib-only.
WorkingDirectory=/root/mammon
ExecStart=/usr/bin/python3 scripts/heartbeat_check.py
# Exit 1 means the heartbeat is stale (alert already POSTed); do not treat the
# oneshot as failed so the timer keeps firing.
SuccessExitStatus=0 1
SyslogIdentifier=mammon-heartbeat
