# scripts/mammon-backup.service
[Unit]
Description=MAMMON SQLite DB backup (oneshot, driven by mammon-backup.timer)
Documentation=https://github.com/kpjmd/Mammon
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
# Stdlib-only online-backup (safe against the live writer) + gzip + prune.
# WorkingDirectory so it reads .env (ALERT_WEBHOOK) and resolves data/ paths.
WorkingDirectory=/opt/mammon
ExecStart=/usr/bin/python3 scripts/backup_db.py --db data/mammon.db --backup-dir data/backups --keep 14
SyslogIdentifier=mammon-backup
# NOTE: off-droplet replication (e.g. `rclone copy data/backups remote:...`) is
# intentionally left to the operator — add it as ExecStartPost once a remote is
# configured, or losing the droplet loses the ledger.
